# install/etc/catkin/profile.d/99_osrf_gear_setup.sh
export GAZEBO_MODEL_PATH=/home/shoky/ariac_ws/install/share/osrf_gear/models:$GAZEBO_MODEL_PATH
